# deploy/systemd/fire-risk-weather-run.service
[Unit]
Description=Fire Risk Hourly Weather Refresh (Open-Meteo)
After=network.target

[Service]
Type=oneshot
User=__SERVICE_USER__
Group=__SERVICE_GROUP__
WorkingDirectory=__REPO_ROOT__
Environment=FIRE_RISK_REPO_ROOT=__REPO_ROOT__
EnvironmentFile=-__REPO_ROOT__/.env
ExecStart=__REPO_ROOT__/scripts/run-weather.sh
NoNewPrivileges=true
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
RestrictSUIDSGID=true
LockPersonality=true
UMask=027
StandardOutput=journal
StandardError=journal
